# s6/php-fpm/finish
#!/bin/sh

kill `cat /var/run/php-fpm.pid`
