# s6/php-fpm/run
#!/bin/execlineb -P

/usr/bin/php-fpm -F
