# s6/php-fpm/run
#!/bin/execlineb -P

/usr/bin/php-fpm --nodaemonize --fpm-config /etc/php5/php-fpm.conf
