# s6/nginx/finish
#!/bin/bash

kill `cat /var/run/nginx.pid`
