# s6/.s6-svscan/finish
#!/bin/sh

for service in /etc/s6/* ; do
   s6-svwait -d $service
done
