# s6/.s6-svscan/finish
#!/bin/sh

for file in /etc/s6/*/finish; do
   $file
done

for service in /etc/s6/* ; do
   s6-svwait -d $service
done
